# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
command=vsftpd

[program:javasandbox]
command=java -jar app.jar
